# systemd/update_index.service
[Unit]
Description=Updates nablaweb DB index

[Service]
User=nablaweb
Group=webapps
Type=oneshot
ExecStart=/etc/websites/nablaweb/update_index.sh