# systemd/update_index.service
[Unit]
Description=Updates nablaweb DB index

[Service]
User=nablaweb
Group=webapps
Type=oneshot
ExecStart=/srv/nablaweb/systemd/update_index.sh
